# var/lib/dpkg/info/qemu-system-common.conffiles
/etc/default/qemu-kvm
/etc/init.d/qemu-kvm
/etc/init/qemu-kvm.conf
/etc/qemu-ifdown
/etc/qemu-ifup
